# templates/savanna-api.conf.erb
description "Savanna API server"
author "Andy Edmonds <andy@edmonds.be>"

start on runlevel [2345]
stop on runlevel [!2345]


chdir /var/run

pre-start script
    mkdir -p /var/run/savanna
    chown savanna:root /var/run/savanna/

    mkdir -p /var/lock/savanna
    chown savanna:root /var/lock/savanna/
end script

exec start-stop-daemon --start --chuid savanna --exec /usr/local/bin/savanna-api -- --config-file=/etc/savanna/savanna.conf
